# sigma-rules/windows/persistence/registry_run_key_persistence.yml
title: Registry Run Key Persistence - Low FP
id: 2b8f4c9a-7f3a-4b1a-9d3e-8a7f2b9d1c11
description: Detects suspicious additions to Windows Run/RunOnce registry keys with filtering to reduce false positives from known software.
status: experimental
author: Sagnik Ray
date: 2026/07/22
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    TargetObject|contains:
      - '\Software\Microsoft\Windows\CurrentVersion\Run'
      - '\Software\Microsoft\Windows\CurrentVersion\RunOnce'
    EventType: 'SetValue'
  filter_known_good:
    Details|contains:
      - 'OneDrive'
      - 'SecurityHealth'
      - 'Windows Defender'
      - 'Teams'
      - 'AdobeGCInvoker'
  condition: selection and not filter_known_good
fields:
  - EventID
  - TargetObject
  - Details
  - Image
falsepositives:
  - Legitimate software installers or updates
level: medium
tags:
  - attack.persistence
  - attack.t1547.001
